# launchd/com.craftcodery.config-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    TEMPLATE REFERENCE ONLY

    This file is NOT deployed directly. The setup.sh script generates
    the actual plist with absolute paths at:
    ~/Library/LaunchAgents/com.craftcodery.config-sync.plist

    This template shows the structure for reference.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.craftcodery.config-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>$HOME/.craftcodery-config/sync.sh --launchd</string>
    </array>

    <!-- Run every hour -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Run on load (after login) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on failure -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>$HOME</string>
        <key>OP_ACCOUNT</key>
        <string>craftcodery.1password.com</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>$HOME</string>

    <key>StandardOutPath</key>
    <string>$HOME/Library/Logs/craftcodery-config-sync-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>$HOME/Library/Logs/craftcodery-config-sync-launchd.log</string>
</dict>
</plist>
